# debian-11/debian-11-sync.service
[Unit]
Description=Manual Debian Mirror Sync Job
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/mirror-sync/debian-11/build-and-sync.sh
