# debian-12/debian-12-sync.service
[Unit]
Description=Manual Debian Mirror Sync Job
After=network.target

[Service]
Type=oneshot
ExecStart=/absolute/path/to/build-and-sync.sh
